# packaging/scripts/preremove.sh
#!/bin/sh
set -eu

if command -v systemctl >/dev/null 2>&1; then
  systemctl stop ech0.service >/dev/null 2>&1 || true
  systemctl disable ech0.service >/dev/null 2>&1 || true
fi
